# initrd/skel/etc/inittab
#
#	$Id$
#

::sysinit:/bin/startup

::once:/bin/installer

tty2::askfirst:/bin/login -f root
tty3::askfirst:/bin/login -f root
tty4::once:/usr/bin/tail -f /log/monitor.log

cons::respawn:/sbin/getty -L 42 console

::ctrlaltdel:/sbin/reboot
