# initrd/skel/etc/inittab
#
#	$Id$
#

::sysinit:/etc/sali/bin/startup

::once:/etc/sali/bin/installer

tty2::askfirst:/bin/login -f root
tty3::askfirst:/bin/login -f root

::ctrlaltdel:/sbin/reboot
