# etc/inittab
#
#	$Id: inittab 439 2011-12-28 14:34:31Z bas $
#

::sysinit:/sbin/sali_startup

::respawn:/sbin/sali_installer

tty2::askfirst:/bin/login -f root
tty3::askfirst:/bin/login -f root
tty4::once:/usr/bin/tail -n1000 -f /var/log/sali.log
tty5::once:/usr/bin/tail -n1000 -f /var/log/sali_console.log

::ctrlaltdel:/sbin/reboot
